# deploy/systemd/orchestrator.service
[Unit]
Description=Smart Outbound Call Orchestrator
After=network.target freeswitch.service redis.service postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agent-flow
ExecStart=/opt/agent-flow/venv/bin/python main.py
Restart=always
RestartSec=5
LimitNOFILE=1000000
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
