# deploy/systemd/qwen-llm.service
[Unit]
Description=Qwen3.5-9B LLM Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/qwen-llm
ExecStart=/opt/qwen-llm/venv/bin/python -m vllm.entrypoints.openai.api_server --model Qwen/Qwen3.5-9B --port 8080
Restart=always
RestartSec=10
LimitNOFILE=1000000
Environment=CUDA_VISIBLE_DEVICES=2

[Install]
WantedBy=multi-user.target
